# openstack/rocky/glance/supervisord.conf
[supervisord]
nodaemon=true

[program:glance-api]
command=glance-api
redirect_stderr=true
